# Disable CPU Boost at boot.service
Service running at /etc/systemd/system/cpufreq-boost-disable.service
[Unit]
Description=Disable CPU Frequency Boost
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo 0 > /sys/devices/system/cpu/cpufreq/boost'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target