# SSD Writes Daily Update.service
Service running at /etc/systemd/system/ssd-daily-writes.service
[Unit]
Description=SSD Daily Write Monitor (SMART + Telegram)
DefaultDependencies=no

# Run VERY early during shutdown
Before=shutdown.target reboot.target halt.target poweroff.target
Conflicts=shutdown.target reboot.target halt.target poweroff.target

# Ensure filesystems are available at boot
After=local-fs.target

[Service]
Type=oneshot

# Capture baseline at boot
ExecStart=/home/sughosha/System_Scripts/System_Health_Monitor/monitors/telegram/ssd_daily_writes_update.sh boot

# Send report at shutdown
ExecStop=/home/sughosha/System_Scripts/System_Health_Monitor/monitors/telegram/ssd_daily_writes_update.sh shutdown

RemainAfterExit=yes

# Absolute safety: never block shutdown
TimeoutStopSec=30
KillMode=process

[Install]
WantedBy=multi-user.target
